# files/etc/uci-defaults/99-yggmesh
#!/bin/sh
# YggMesh Node — Zero-Touch First Boot Configuration
# This script runs once on first boot via uci-defaults mechanism.
# It configures: 802.11s WiFi mesh, dual-band WiFi (private AP), firewall zones, Yggdrasil overlay.

# ============================================================================
# 1. IDENTITY GENERATION
# ============================================================================

# Client AP SSID for end-user devices
PRIVATE_SSID="YggMesh"
# 802.11s mesh identifier (must match on all nodes)
MESH_ID="yggmesh/mesh"
# SAE key for 802.11s mesh peering (must match on all nodes)
MESH_KEY="qJ7tN2vL8pR4xKcM"
# Initial root password set on first boot (user should change via LuCI or passwd)
DEFAULT_ROOT_PASSWORD="yggmesh"
# Yggdrasil TUN interface name
YGG_IFNAME="ygg0"
# Fixed router IPv4 address for local WiFi
PRIV_IP="10.0.0.1"
PRIV_SUBNET="10.0.0.0"
# Yggdrasil fixed peering port
YGG_PORT="17000"
# Yggdrasil multicast interface regex
YGG_PEER_REGEX=".+"
# Alfis DNS resolver (Yggdrasil global DNS)
ALFIS_DNS="324:71e:281a:9ed3::53"

if [ -f /etc/yggmesh/keys ] && grep -q "^NODE_MAC=" /etc/yggmesh/keys; then
    # OTA update — /etc/yggmesh/keys preserved by sysupgrade.conf
    # Load existing identity to keep hostname/MAC stable across firmware updates
    . /etc/yggmesh/keys
    BASE_MAC="$NODE_MAC"
    NODE_SUFFIX=$(echo "$BASE_MAC" | awk -F: '{print toupper($5$6)}')
    HOSTNAME="YggMesh-${NODE_SUFFIX}"
    logger -t yggmesh "OTA: preserving identity ${HOSTNAME} (MAC ${BASE_MAC})"
else
    # Fresh flash — detect MAC from hardware
    # Get hardware MAC with retry (some devices initialize NICs slowly at first boot)
    BASE_MAC=""
    for attempt in 1 2 3 4 5; do
        for iface in br-lan eth0 eth1 eth2; do
            mac=$(cat /sys/class/net/$iface/address 2>/dev/null)
            if [ -n "$mac" ] && [ "$mac" != "00:00:00:00:00:00" ]; then
                BASE_MAC="$mac"
                break 2
            fi
        done
        [ "$attempt" -lt 5 ] && sleep 1
    done
    # Fallback: random locally-administered MAC (02:xx:xx:xx:xx:xx)
    if [ -z "$BASE_MAC" ] || [ "$BASE_MAC" = "00:00:00:00:00:00" ]; then
        BASE_MAC=$(hexdump -n 5 -e '/1 ":%02x"' /dev/urandom | sed 's/^:/02:/')
    fi

    # Last 4 hex digits of MAC for unique suffix
    NODE_SUFFIX=$(echo "$BASE_MAC" | awk -F: '{print toupper($5$6)}')
    HOSTNAME="YggMesh-${NODE_SUFFIX}"
fi

# ============================================================================
# 3. KEY GENERATION
# ============================================================================



mkdir -p /etc/yggmesh
cat > /etc/yggmesh/keys <<KEYS_EOF
# YggMesh Node Keys — Generated at first boot
# Node: ${HOSTNAME}
# Date: $(date -u +%Y-%m-%dT%H:%M:%SZ)
MESH_SAE_KEY=${MESH_KEY}
PRIVATE_SSID=${PRIVATE_SSID}
PRIVATE_SUBNET=${PRIV_SUBNET}/24
NODE_MAC=${BASE_MAC}
KEYS_EOF
chmod 600 /etc/yggmesh/keys

# Set default root password (user should change via LuCI or passwd)
echo -e "${DEFAULT_ROOT_PASSWORD}\n${DEFAULT_ROOT_PASSWORD}" | passwd root >/dev/null 2>&1

# ============================================================================
# 4. NETWORK CONFIGURATION
# ============================================================================

# --- Port mapping (written by build.sh per device) ---
PORT_MAP=$(cat /etc/yggmesh/port_map 2>/dev/null || echo "eth0:lan eth1:wan")
if [ "$PORT_MAP" = "dsa" ]; then
    # DSA switch (ax53u) — individual ports
    LAN_PORTS="lan1 lan2 lan3 lan4"
    WAN_DEV="wan"
else
    # Parse port map: "eth0:wan eth1:lan" → extract LAN and WAN
    LAN_PORTS=""
    WAN_DEV=""
    for mapping in $PORT_MAP; do
        port="${mapping%%:*}"
        role="${mapping##*:}"
        if [ "$role" = "lan" ]; then
            LAN_PORTS="$port"
        elif [ "$role" = "wan" ]; then
            WAN_DEV="$port"
        fi
    done
    [ -z "$LAN_PORTS" ] && LAN_PORTS="eth0"
    [ -z "$WAN_DEV" ] && WAN_DEV="eth1"
fi

# Delete default OpenWrt br-lan/lan (conflicts with our br-private)
uci -q delete network.lan 2>/dev/null || true
# Find and delete br-lan device
for i in 0 1 2 3 4; do
    name=$(uci -q get "network.@device[$i].name" 2>/dev/null)
    if [ "$name" = "br-lan" ]; then
        uci delete "network.@device[$i]"
        break
    fi
done

uci batch <<-NET_EOF
# --- loopback: Standard local routing interface ---
set network.loopback=interface
set network.loopback.device='lo'
set network.loopback.proto='static'
set network.loopback.ipaddr='127.0.0.1'
set network.loopback.netmask='255.0.0.0'

# --- private_dev: Physical bridge grouping LAN ethernet ports and the Client AP Wi-Fi ---
set network.private_dev=device
set network.private_dev.type='bridge'
set network.private_dev.name='br-private'
set network.private_dev.empty='1'
set network.private_dev.igmp_snooping='0'
delete network.private_dev.ports

# --- mesh: Layer 2 interface for raw 802.11s radios. proto='unmanaged' ensures it gets an IPv6 Link-Local address ---
set network.mesh=interface
set network.mesh.proto='unmanaged'

# --- private: Layer 3 logical interface for br-private, providing the local router IP to clients ---
set network.private=interface
set network.private.device='br-private'
set network.private.proto='static'
set network.private.ipaddr='${PRIV_IP}'
set network.private.netmask='255.255.255.0'

# --- wan: Upstream internet port (DHCP client) to connect to global Yggdrasil peers ---
set network.wan=interface
set network.wan.device='${WAN_DEV}'
set network.wan.proto='dhcp'
NET_EOF

# Add LAN ports to private bridge (DSA: lan1-4, non-DSA: eth0)
for port in $LAN_PORTS; do
    uci add_list network.private_dev.ports="$port"
done

uci commit network

# ============================================================================
# 5. WIRELESS CONFIGURATION
# ============================================================================

# Detect radios dynamically
RADIOS=$(uci show wireless | grep '=wifi-device' | cut -d'=' -f1 | cut -d'.' -f2)

# Determine which radio is 2.4GHz and which is 5GHz
RADIO_2G=""
RADIO_5G=""

for radio in $RADIOS; do
    band=$(uci -q get "wireless.${radio}.band" 2>/dev/null || echo "")
    htmode=$(uci -q get "wireless.${radio}.htmode" 2>/dev/null || echo "")

    if [ "$band" = "2g" ]; then
        RADIO_2G="$radio"
    elif [ "$band" = "5g" ]; then
        RADIO_5G="$radio"
    else
        # Fallback: detect from htmode or channel
        channel=$(uci -q get "wireless.${radio}.channel" 2>/dev/null || echo "0")
        if [ "$channel" -le 14 ] 2>/dev/null; then
            RADIO_2G="$radio"
        else
            RADIO_5G="$radio"
        fi
    fi
done

# Fallback: if only one radio detected, use it for both
DUAL_BAND=true
if [ -z "$RADIO_2G" ] && [ -n "$RADIO_5G" ]; then
    RADIO_2G="$RADIO_5G"
    DUAL_BAND=false
fi
if [ -z "$RADIO_5G" ] && [ -n "$RADIO_2G" ]; then
    RADIO_5G="$RADIO_2G"
    DUAL_BAND=false
fi

# Remove all existing wifi-iface sections
while uci -q delete wireless.@wifi-iface[0]; do :; done

# Enable radios and set fixed channels (all nodes must share same channel for mesh peering)
for radio in $RADIOS; do
    uci -q delete "wireless.${radio}.disabled" 2>/dev/null || true
done

# Fixed channels — all mesh nodes must be on the same channel per band
# 2.4GHz: channel 1 (DFS-free, least congested of 1/6/11)
# 5GHz: channel 36 (DFS-free, universal chipset support)
if [ -n "$RADIO_2G" ]; then
    uci set "wireless.${RADIO_2G}.channel=1"
fi
if [ -n "$RADIO_5G" ] && [ "$DUAL_BAND" = "true" ]; then
    uci set "wireless.${RADIO_5G}.channel=36"
fi

# --- 2.4GHz radio: always mesh backhaul; AP only if single-band ---
if [ -n "$RADIO_2G" ]; then
    # Mesh interface on 2.4GHz (always — dedicated backhaul on all devices)
    uci batch <<-WIFI_2G_MESH
set wireless.mesh_2g=wifi-iface
set wireless.mesh_2g.device='${RADIO_2G}'
set wireless.mesh_2g.mode='mesh'
set wireless.mesh_2g.mesh_id='${MESH_ID}'
set wireless.mesh_2g.mesh_fwding='0'
set wireless.mesh_2g.encryption='none'
set wireless.mesh_2g.network='mesh'
WIFI_2G_MESH

    # Private AP on 2.4GHz — single-band devices only
    # (dual-band devices dedicate 2.4GHz to mesh and use 5GHz for client AP)
    if [ "$DUAL_BAND" = "false" ]; then
        uci batch <<-WIFI_2G_PRIV
set wireless.private_2g=wifi-iface
set wireless.private_2g.device='${RADIO_2G}'
set wireless.private_2g.mode='ap'
set wireless.private_2g.ssid='${PRIVATE_SSID}'
set wireless.private_2g.encryption='none'
set wireless.private_2g.network='private'
WIFI_2G_PRIV
    fi
fi

# --- 5GHz radio: AP only (dual-band devices) — dedicated client radio ---
if [ -n "$RADIO_5G" ] && [ "$DUAL_BAND" = "true" ]; then
    uci batch <<-WIFI_5G_PRIV
set wireless.private_5g=wifi-iface
set wireless.private_5g.device='${RADIO_5G}'
set wireless.private_5g.mode='ap'
set wireless.private_5g.ssid='${PRIVATE_SSID}'
set wireless.private_5g.encryption='none'
set wireless.private_5g.network='private'
WIFI_5G_PRIV
fi

uci commit wireless

# ============================================================================
# 6. FIREWALL CONFIGURATION
# ============================================================================

# Reset firewall to clean state
while uci -q delete firewall.@zone[0]; do :; done
while uci -q delete firewall.@forwarding[0]; do :; done
while uci -q delete firewall.@rule[0]; do :; done
while uci -q delete firewall.@redirect[0]; do :; done

uci batch <<-FW_EOF
# --- Zone: lan (private network) ---
add firewall zone
set firewall.@zone[-1].name='lan'
set firewall.@zone[-1].input='ACCEPT'
set firewall.@zone[-1].output='ACCEPT'
set firewall.@zone[-1].forward='ACCEPT'
add_list firewall.@zone[-1].network='private'

# --- Zone: mesh (raw 802.11s links for Yggdrasil) ---
add firewall zone
set firewall.@zone[-1].name='mesh'
set firewall.@zone[-1].input='REJECT'
set firewall.@zone[-1].output='ACCEPT'
set firewall.@zone[-1].forward='REJECT'
add_list firewall.@zone[-1].network='mesh'

# Allow incoming Yggdrasil IPv6 traffic on mesh
add firewall rule
set firewall.@rule[-1].name='Allow-Yggdrasil-Mesh'
set firewall.@rule[-1].src='mesh'
set firewall.@rule[-1].proto='tcp udp'
set firewall.@rule[-1].dest_port='${YGG_PORT}'
set firewall.@rule[-1].family='ipv6'
set firewall.@rule[-1].target='ACCEPT'

# Allow IPv6 Link-Local routing (fe80::) unconditionally on mesh for Yggdrasil beacons
add firewall rule
set firewall.@rule[-1].name='Allow-LinkLocal-Mesh'
set firewall.@rule[-1].src='mesh'
set firewall.@rule[-1].src_ip='fe80::/10'
set firewall.@rule[-1].family='ipv6'
set firewall.@rule[-1].target='ACCEPT'

# Allow IPv6 Multicast (ff02::) unconditionally on mesh for Yggdrasil discovery
add firewall rule
set firewall.@rule[-1].name='Allow-Multicast-Mesh'
set firewall.@rule[-1].src='mesh'
set firewall.@rule[-1].dest_ip='ff00::/8'
set firewall.@rule[-1].family='ipv6'
set firewall.@rule[-1].target='ACCEPT'

# Allow ICMPv6 on mesh (needed for IPv6 Neighbor Discovery and SLAAC)
add firewall rule
set firewall.@rule[-1].name='Allow-ICMPv6-Mesh'
set firewall.@rule[-1].src='mesh'
set firewall.@rule[-1].proto='icmp'
set firewall.@rule[-1].family='ipv6'
set firewall.@rule[-1].target='ACCEPT'

# --- Zone: wan ---
add firewall zone
set firewall.@zone[-1].name='wan'
set firewall.@zone[-1].input='REJECT'
set firewall.@zone[-1].output='ACCEPT'
set firewall.@zone[-1].forward='REJECT'
set firewall.@zone[-1].masq='1'
set firewall.@zone[-1].mtu_fix='1'
add_list firewall.@zone[-1].network='wan'

# --- Forwarding: lan → wan ---
add firewall forwarding
set firewall.@forwarding[-1].src='lan'
set firewall.@forwarding[-1].dest='wan'

# --- Standard WAN input rules ---
add firewall rule
set firewall.@rule[-1].name='Allow-DHCP-Renew'
set firewall.@rule[-1].src='wan'
set firewall.@rule[-1].proto='udp'
set firewall.@rule[-1].dest_port='68'
set firewall.@rule[-1].target='ACCEPT'
set firewall.@rule[-1].family='ipv4'

add firewall rule
set firewall.@rule[-1].name='Allow-Ping'
set firewall.@rule[-1].src='wan'
set firewall.@rule[-1].proto='icmp'
set firewall.@rule[-1].icmp_type='echo-request'
set firewall.@rule[-1].target='ACCEPT'
FW_EOF

uci commit firewall

# ============================================================================
# 7. DHCP / DNS CONFIGURATION
# ============================================================================

# Remove existing DHCP pools that may conflict
uci -q delete dhcp.lan 2>/dev/null || true

uci batch <<-DHCP_EOF
set dhcp.private=dhcp
set dhcp.private.interface='private'
set dhcp.private.start='100'
set dhcp.private.limit='150'
set dhcp.private.leasetime='12h'
# Advertise router's own IPv4 for DNS (which forwards to Alfis)
# Note: DHCPv4 option 6 cannot be an IPv6 address
add_list dhcp.private.dhcp_option='6,${PRIV_IP}'

# Configure Dnsmasq to use Alfis, alongside reliable public DNS
# This bypasses potentially broken/hijacked ISP DNS so public peers always resolve
set dhcp.@dnsmasq[0].noresolv='1'
set dhcp.@dnsmasq[0].localuse='1'
add_list dhcp.@dnsmasq[0].server='${ALFIS_DNS}'
add_list dhcp.@dnsmasq[0].server='1.1.1.1'
add_list dhcp.@dnsmasq[0].server='8.8.8.8'
DHCP_EOF

uci commit dhcp

# ============================================================================
# 8. SYSTEM SETTINGS
# ============================================================================

uci batch <<-SYS_EOF
set system.@system[0].hostname='${HOSTNAME}'
set system.@system[0].timezone='UTC0'
set system.@system[0].zonename='UTC'
set system.@system[0].log_proto='udp'
set system.@system[0].conloglevel='7'
set system.@system[0].cronloglevel='9'
SYS_EOF

# NTP: use yggmesh-ng via Yggdrasil as primary server (available even without WAN)
uci -q delete system.ntp.server
uci add_list system.ntp.server='<YGG_ADDR_MAIN>'
uci add_list system.ntp.server='0.openwrt.pool.ntp.org'
uci add_list system.ntp.server='1.openwrt.pool.ntp.org'
uci add_list system.ntp.server='2.openwrt.pool.ntp.org'
uci set system.ntp.enabled='1'

uci commit system


# =========================================================================
# 9. YGGDRASIL OVERLAY NETWORK
# =========================================================================

# Private key: generated once, stored in /etc/yggmesh/keys (preserved across OTA by sysupgrade.conf)
# The native OpenWrt UCI is configured so the LuCI app can manage Yggdrasil peers.
if grep -q "^YGG_PRIVKEY=" /etc/yggmesh/keys 2>/dev/null; then
    YGG_PRIVKEY=$(grep "^YGG_PRIVKEY=" /etc/yggmesh/keys | cut -d= -f2-)
else
    # Fresh install: generate a new key pair and persist the private key
    YGG_PRIVKEY=$(yggdrasil -genconf 2>/dev/null | grep 'PrivateKey' | awk -F': ' '{print $2}')
    echo "YGG_PRIVKEY=${YGG_PRIVKEY}" >> /etc/yggmesh/keys
fi

# Configure Yggdrasil interface natively via UCI
uci batch <<-YGG_EOF
set network.yggdrasil=interface
set network.yggdrasil.proto='yggdrasil'
set network.yggdrasil.private_key="${YGG_PRIVKEY}"
set network.yggdrasil.node_info_privacy='0'

add network yggdrasil_yggdrasil_interface
add_list network.@yggdrasil_yggdrasil_interface[-1].interface='${YGG_PEER_REGEX}'
set network.@yggdrasil_yggdrasil_interface[-1].beacon='1'
set network.@yggdrasil_yggdrasil_interface[-1].listen='1'
set network.@yggdrasil_yggdrasil_interface[-1].port='${YGG_PORT}'
YGG_EOF
uci commit network

# Yggdrasil firewall zone
uci batch <<-YGG_FW_EOF
add firewall zone
set firewall.@zone[-1].name='yggdrasil'
set firewall.@zone[-1].input='REJECT'
set firewall.@zone[-1].output='ACCEPT'
set firewall.@zone[-1].forward='REJECT'
add_list firewall.@zone[-1].network='yggdrasil'

# Allow ICMPv6 on Yggdrasil (Required)
add firewall rule
set firewall.@rule[-1].name='Allow ICMPv6 Ygg'
set firewall.@rule[-1].src='yggdrasil'
set firewall.@rule[-1].proto='icmp'
set firewall.@rule[-1].family='ipv6'
set firewall.@rule[-1].target='ACCEPT'

# Forwardings: LAN → Yggdrasil (outbound: private clients can reach Yggdrasil services)
add firewall forwarding
set firewall.@forwarding[-1].src='lan'
set firewall.@forwarding[-1].dest='yggdrasil'
YGG_FW_EOF

uci commit firewall

# Enable yggdrasil service
/etc/init.d/yggdrasil enable 2>/dev/null || true

# Temporary config for address generation
YGG_TMP=$(mktemp)
echo "{ PrivateKey: \"${YGG_PRIVKEY}\" }" > "$YGG_TMP"

# Save yggdrasil address to node keys file
YGG_ADDR=$(yggdrasil -address -useconffile "$YGG_TMP" 2>/dev/null || echo "unknown")
echo "YGGDRASIL_ADDRESS=${YGG_ADDR}" >> /etc/yggmesh/keys

# ============================================================================
# 10. YGGDRASIL IPv6
#     Private: 300::/64 SLAAC — devices reachable via Yggdrasil
# ============================================================================

# Extract Yggdrasil 300::/64 subnet
YGG_SUBNET=$(yggdrasil -subnet -useconffile "$YGG_TMP" 2>/dev/null)
rm -f "$YGG_TMP"

# YGG_SUBNET format: 300:xxxx:xxxx:xxxx::/64
YGG_SUBNET_PREFIX="${YGG_SUBNET%%::*}"  # 300:xxxx:xxxx:xxxx

if [ -n "$YGG_SUBNET_PREFIX" ]; then
    # Assign 300::/64 to PRIVATE interface
    # Private clients get globally-routable Yggdrasil addresses via SLAAC
    uci set network.private.ip6addr="${YGG_SUBNET_PREFIX}::1/64"
    uci commit network

    # Configure RA/SLAAC for private (odhcpd)
    uci batch <<-PRIV6_EOF
set dhcp.private.ra='server'
set dhcp.private.dhcpv6='disabled'
set dhcp.private.ra_management='0'
set dhcp.private.ra_default='1'
set dhcp.private.ra_slaac='1'
add_list dhcp.private.dns='${ALFIS_DNS}'
PRIV6_EOF
    uci commit dhcp

    # Save subnet to keys file for reference
    echo "YGGDRASIL_SUBNET=${YGG_SUBNET}" >> /etc/yggmesh/keys

    logger -t yggmesh "Private IPv6: ${YGG_SUBNET_PREFIX}::1/64 (SLAAC)"
else
    logger -t yggmesh "Warning: Could not extract Yggdrasil subnet"
fi

# ============================================================================
# 11. FINAL
# ============================================================================

# Log completion
logger -t yggmesh "First boot configuration complete"
logger -t yggmesh "Hostname: ${HOSTNAME}"
logger -t yggmesh "Private: ${PRIVATE_SSID} @ ${PRIV_IP}/24"
logger -t yggmesh "Mesh ID: ${MESH_ID}"
logger -t yggmesh "Yggdrasil: ${YGG_ADDR} (multicast discovery on br-private)"
logger -t yggmesh "Private IPv6: Yggdrasil 300::/64 SLAAC"

exit 0
